# ansible_setup/roles/monitor/files/monitor.conf
[program:monitor]
directory=/opt/
command=/opt/monitor.sh
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/monitor.log
logfile_maxbytes=5M
stopasgroup=true
stopsignal=QUIT