# ansible_setup/roles/manipulatord/files/supervisord.conf
[program:siit-manipulatord]
directory=/home/ubuntu/siit-manipulatord
command=/home/ubuntu/siit-manipulatord/start.sh
numprocs=1
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/log/supervisor/siit-manipulatord.log
logfile_maxbytes=5M
stopasgroup=true
stopsignal=QUIT